# recorder-manager/rootfs/etc/services.d/recorder-manager/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Recorder Manager service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting Recorder Manager..."

exec python3 /usr/src/app/server.py
